# .config/systemd/user/emacs.service.d/override.conf
[Service]
ExecStart=
ExecStart=/usr/bin/emacs --init-directory ~/.emacs.d --fg-daemon
